# swift/Tests/ADOpsTests/Fixtures/expected/mt.everest.apple-docs.cloudflared.web.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>mt.everest.apple-docs.cloudflared.web</string>
    <key>UserName</key>
    <string>everest</string>
    <key>GroupName</key>
    <string>staff</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/everest/bin/cloudflared</string>
        <string>tunnel</string>
        <string>--no-autoupdate</string>
        <string>--config</string>
        <string>/Users/everest/Developer/apple-docs/ops/cloudflared/config.yml</string>
        <string>run</string>
        <string>apple-docs</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/everest/Developer/apple-docs/ops</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>ThrottleInterval</key>
    <integer>10</integer>
    <key>ProcessType</key>
    <string>Background</string>
    <key>StandardOutPath</key>
    <string>/Users/everest/Developer/apple-docs/ops/logs/cloudflared-web.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/everest/Developer/apple-docs/ops/logs/cloudflared-web.err.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>HOME</key>
        <string>/Users/everest</string>
        <key>TUNNEL_ORIGIN_CERT</key>
        <string>/Users/everest/.cloudflared/cert.pem</string>
    </dict>
</dict>
</plist>
